# server/supervisor.conf
[program:django_wx_tornado-8081]
command=/usr/local/python/bin/python index.py
#command=python index.py
stderr_logfile=tornado-stderr.log
stdout_logfile=tornado-stdout.log

[supervisord]
